# units/collect_metrics.service
[Unit]
Description=Collect Metrics Service

[Service]
ExecStart=/bin/bash -c '/etc/systemd/user/collect_metrics.sh'
