# units/panamax.service
[Unit]
Description=Panamax local
After=docker.service
Requires=docker.service

[Service]
ExecStart=/etc/systemd/user/start_svc.sh

[Install]
WantedBy=multi-user.target
